# deploy/com.caddy.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label must match filename (without .plist) -->
    <key>Label</key>
    <string>com.caddy.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/caddy</string>
        <string>run</string>
        <string>--config</string>
        <!-- Update this path to match your actual repo location -->
        <string>/Users/abraxas/Personal/claudio/deploy/Caddyfile</string>
    </array>

    <!-- Start automatically on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/caddy.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/caddy.err</string>
</dict>
</plist>
